# configs/com.granola-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.granola-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/granola-sync</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>~/.config/granola-sync/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>~/.config/granola-sync/stderr.log</string>

    <key>WorkingDirectory</key>
    <string>~</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
